# markdown-wiki.service
[Unit]
Description=Create the markdown-wiki

[Service]
ExecStart=/usr/bin/mdwiki.sh /share/data/wiki/ /var/www/html/wiki/
User=root
