# markdown-wiki.service
[Unit]
Description=Create the markdown-wiki

[Service]
ExecStart=/opt/clang/markdown-wiki/mdwiki.sh /share/data/wiki/ /var/www/html/wiki/
User=root
